# scripts/huawei/systemd/patch-lifecycle.service.template
[Unit]
Description=Patch Lifecycle Compose Service
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__APP_DIR__
ExecStart=/usr/bin/env bash -lc "cd __APP_DIR__ && docker compose -f __COMPOSE_FILE__ up -d"
ExecStop=/usr/bin/env bash -lc "cd __APP_DIR__ && docker compose -f __COMPOSE_FILE__ down"
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
